# infra/systemd/polymarket-crawler.service
[Unit]
Description=Polymarket Wallet Universe Crawler (depth-tier maintenance + backfill)
After=network-online.target postgresql.service redis-server.service polymarket-onchain.service
Wants=network-online.target

[Service]
Type=simple
User=polymarket
Group=polymarket
WorkingDirectory=/opt/polymarket-bot/
EnvironmentFile=/opt/polymarket-bot/.env
# Wave 2 adds src/crawler/main.py as the daemon entrypoint.
ExecStart=/opt/polymarket-bot/.venv/bin/python -m src.crawler.main
Restart=always
RestartSec=5s
StandardOutput=journal
StandardError=journal

# Round 6 § 3.5 memory budget.
MemoryMax=200M

[Install]
WantedBy=multi-user.target
