# infra/systemd/polymarket-social.service
[Unit]
Description=Polymarket Bot — Social signal daemon (Round 12)
Documentation=https://github.com/your-org/polymarket-bot/blob/main/docs/ROUND_12_SOCIAL_AND_CROSS_MARKET.md
After=network-online.target postgresql.service redis-server.service
Wants=network-online.target

[Service]
Type=simple
User=polymarket
Group=polymarket
WorkingDirectory=/opt/polymarket-bot/
EnvironmentFile=/opt/polymarket-bot/.env
ExecStart=/opt/polymarket-bot/.venv/bin/python -m src.social
Restart=always
RestartSec=5s

# Per Round 12 § 2 — the social daemon hosts X firehose + NLP
# classifier + Telegram/Discord listeners. Spec target ~250 MB; we
# allow 300 MB headroom (operator can step down post-stabilisation).
MemoryMax=300M

StandardOutput=journal
StandardError=journal

SyslogIdentifier=polymarket-social

[Install]
WantedBy=multi-user.target
